# deploy/systemd/ai-employee-twitter.service
[Unit]
Description=AI Employee — Twitter/X Watcher (Gold)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__AI_USER__
Group=__AI_USER__
WorkingDirectory=__VAULT_PATH__
EnvironmentFile=-__VAULT_PATH__/.env
ExecStart=__VAULT_PATH__/.venv/bin/python3 watchers/twitter_watcher.py \
    --vault-path __VAULT_PATH__
Restart=on-failure
RestartSec=30
StartLimitIntervalSec=600
StartLimitBurst=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-twitter

[Install]
WantedBy=multi-user.target
