# deploy/systemd/ai-employee-briefing.service
[Unit]
Description=AI Employee — Daily Briefing Generator (Silver)
After=network.target

[Service]
Type=oneshot
User=__AI_USER__
Group=__AI_USER__
WorkingDirectory=__VAULT_PATH__
EnvironmentFile=-__VAULT_PATH__/.env
ExecStart=__VAULT_PATH__/.venv/bin/python3 watchers/daily_briefing_scheduler.py \
    --vault-path __VAULT_PATH__
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-briefing
