# deploy/systemd/ai-employee-cloud-orchestrator.service
[Unit]
Description=AI Employee — Cloud Orchestrator (Platinum)
After=network-online.target vault-sync.service
Wants=network-online.target
BindsTo=vault-sync.service

[Service]
Type=simple
User=__AI_USER__
Group=__AI_USER__
WorkingDirectory=__VAULT_PATH__
EnvironmentFile=-__VAULT_PATH__/.env
Environment=AGENT_ID=cloud
ExecStart=__VAULT_PATH__/.venv/bin/python3 watchers/cloud_orchestrator.py \
    --vault-path __VAULT_PATH__ \
    --agent-id cloud
Restart=on-failure
RestartSec=30
StartLimitIntervalSec=600
StartLimitBurst=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-cloud-orchestrator

[Install]
WantedBy=multi-user.target
