# deploy/systemd/ai-employee-briefing.timer
[Unit]
Description=AI Employee — Daily Briefing Timer (8:00 AM daily)
Requires=ai-employee-briefing.service

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
