# deploy/systemd/ai-employee-filesystem.service
[Unit]
Description=AI Employee — Filesystem Watcher (Bronze)
After=network.target
Wants=network.target

[Service]
Type=simple
User=__AI_USER__
Group=__AI_USER__
WorkingDirectory=__VAULT_PATH__
EnvironmentFile=-__VAULT_PATH__/.env
ExecStart=__VAULT_PATH__/.venv/bin/python3 watchers/filesystem_watcher.py \
    --vault-path __VAULT_PATH__
Restart=on-failure
RestartSec=15
StartLimitIntervalSec=300
StartLimitBurst=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-filesystem

[Install]
WantedBy=multi-user.target
