# deploy/systemd/ai-employee-health.service
[Unit]
Description=AI Employee — Health Monitor (Gold)
After=multi-user.target
Wants=ai-employee-filesystem.service ai-employee-gmail.service

[Service]
Type=simple
User=__AI_USER__
Group=__AI_USER__
WorkingDirectory=__VAULT_PATH__
EnvironmentFile=-__VAULT_PATH__/.env
ExecStart=__VAULT_PATH__/.venv/bin/python3 watchers/health_monitor.py \
    --vault-path __VAULT_PATH__
Restart=always
RestartSec=60
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-health

[Install]
WantedBy=multi-user.target
