# deploy/systemd/vault-sync.timer
[Unit]
Description=AI Employee — Vault Sync Timer (every 5 minutes)
Requires=vault-sync.service

[Timer]
OnBootSec=60
OnUnitActiveSec=5min
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
